# systemd/zenbook-thermal.service
[Unit]
Description=Zenbook Duo Thermal Monitor (auto fan profile)
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/thermal-monitor.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
