# systemd/zenbook-light-monitor.service
[Unit]
Description=Zenbook Duo Keyboard Backlight (ambient light + inactivity)
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kb-backlight-unified.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
