# systemd/zenbook-auto-display.service
[Unit]
Description=Zenbook Duo Auto Display
After=graphical.target

[Service]
Type=forking
ExecStart=/home/carlosh/asus_UX8406MA/zenbook-duo-linux/scripts/auto-display.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target