# systemd/zenbook-nightlight.service
[Unit]
Description=Zenbook Duo Night Light (Redshift)
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/nightlight.sh auto
Restart=on-failure
RestartSec=60

[Install]
WantedBy=graphical.target
